# dockerfiles/standalone/lab41/spark-master/config/service/ssh/run
#!/bin/sh

# redirect stderr->stdout
exec 2>&1

# start ssh daemon
exec /usr/sbin/sshd -D -e
